# deploy/cryptoapp.service
# Systemd unit file for CryptoApp
# Instructions:
# 1. Copy to /etc/systemd/system/cryptoapp.service
# 2. Replace YOUR_USERNAME with your Azure VM username (e.g., azureuser)
# 3. Replace /path/to/CryptoApp with your actual deployment path (e.g., /home/azureuser/CryptoApp)
# 4. Run: sudo systemctl daemon-reload
# 5. Run: sudo systemctl enable cryptoapp
# 6. Run: sudo systemctl start cryptoapp

[Unit]
Description=CryptoApp Flask (Gunicorn)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/CryptoApp
Environment=PYTHONUNBUFFERED=1
# Using uv to ensure deps from pyproject/uv.lock; remove `uv sync` if you manage venv separately
ExecStart=/usr/bin/env bash -lc 'cd /path/to/CryptoApp && uv sync && uv run gunicorn -c gunicorn.conf.py wsgi:app'
Restart=always
RestartSec=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
